# templates/nsqd.service.erb
[Unit]
Description=nsqd
After=network.target

[Service]
ExecStart=<%= scope.lookupvar('nsq::bin_dir')%>/bin/nsqd -config <%= scope.lookupvar('nsq::conf_dir') %>/nsqd.conf
Type=simple
PIDFile=/var/run/nsqd.pid

[Install]
WantedBy=default.target
